# 集群/conf/tomcat/tomcat54080.ini
[program:tomcat54080]
command=/opt/SoftwareManager/tomcat-54080/bin/catalina.sh run
autostart=true
startsecs=40
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups = 20
stdout_logfile=/opt/SoftwareManager/tomcat-54080/logs/catalina.out
stopasgroup=false
killasgroup=false
